# deploy/tg-ingest.service
[Unit]
Description=Telegram -> D1 ingest (off-Cloudflare Telethon fetcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Secrets/config (chmod 600). See tg-ingest.env.example.
EnvironmentFile=/etc/tg-crm-ingest.env
# EDIT these three to match your install location + venv:
User=tgingest
WorkingDirectory=/opt/tg-crm-system
ExecStart=/opt/tg-crm-system/.venv/bin/python /opt/tg-crm-system/scripts/ingest_local.py
# Don't let one slow/hung run pile up or run forever.
TimeoutStartSec=180
# Hardening (safe for a network client that only reads files + makes egress).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
